# daemon/daemon
#!/bin/sh
### BEGIN INIT INFO
# Provides:          RPi-Monitor
# Required-Start:    $remote_fs $syslog $network 
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: This file starts the Loxberry RPi-Monitor update
# Description:       This file starts the Loxberry RPi-Monitor update.
### END INIT INFO

. /lib/lsb/init-functions

### END INIT INFO

<<<<<<< HEAD
# Version 1.1
# 15.09.2016 15:14:00

if [ -f REPLACEBYPATHTODATA/doinstall ]
	then
    log_action_msg "Fresh installation. Installing RPiMonitor package"
    dpkg -i REPLACEBYPATHTODATA/rpimonitor_latest_all.deb
    cp REPLACEBYPATHTODATA/data.conf /etc/rpimonitor/
    cp REPLACEBYPATHTODATA/version.conf /etc/rpimonitor/template/
    rm REPLACEBYPATHTODATA/doinstall
	fi
=======
# Version 1.3
# 15.09.2016 00:02:36

# Prepare System Monitoring Tools
if [ -e /etc/apt/sources.list.d/rpimonitor.list ]
then
	log_action_begin_msg "/etc/apt/sources.list.d/rpimonitor.list seems available, skipping download"
else
	log_action_begin_msg "/etc/apt/sources.list.d/rpimonitor.list not found, try to download"
	wget -T 2 http://goo.gl/vewCLL -O /etc/apt/sources.list.d/rpimonitor.list
	apt-key adv --recv-keys --keyserver keyserver.ubuntu.com 2C0D3C0F 
fi

# If System Monitoring Tools installed, enable autoupdate when reboot
if [ -e /etc/init.d/rpimonitor  ]
then
   log_action_msg "Error running '/etc/init.d/rpimonitor update'"
	/etc/init.d/rpimonitor install_auto_package_status_update 
fi

exit 0
>>>>>>> a22a97dfe474aa1ec6d7f53a52f7a50629d42530
